# brain/services/nox-poller.service
[Unit]
Description=Nox Body Poller (polls PiDog for voice/sensor events)
After=network.target openclaw-gateway.service

[Service]
Type=simple
User=piclawbot
WorkingDirectory=/home/piclawbot/clawd/scripts
ExecStart=/usr/bin/python3 -u /home/piclawbot/clawd/scripts/nox_body_poller.py
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/piclawbot/.openclaw/.env
Environment=PIDOG_HOST=pidog.local
Environment=PIDOG_BRIDGE_PORT=8888
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
